# docker/joomla/build-files/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=service php5-fpm start

[program:apache2]
command=service apache2 start
